# datos/procesados/launchd/com.masterescala.comunidad-vip-watcher.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.masterescala.comunidad-vip-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__INSTALL_DIR__/pipeline_archivo_nuevo.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>

    <!-- ¡La magia! macOS dispara este launchd cuando algo cambia en esta carpeta. -->
    <key>WatchPaths</key>
    <array>
        <string>__ORIGINALES_DIR__</string>
    </array>

    <!-- Throttle: no disparar más de una vez cada 60 seg, aunque haya múltiples eventos -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>__INSTALL_DIR__/logs/launchd-watcher.out.log</string>

    <key>StandardErrorPath</key>
    <string>__INSTALL_DIR__/logs/launchd-watcher.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
